# backend/systemd/cmp-backend.service
[Unit]
Description=CMP Backend Service
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/cmp/backend
Environment=NODE_ENV=production
EnvironmentFile=/srv/cmp/backend/.env
ExecStart=/usr/bin/env node index.js
Restart=always
RestartSec=3s
StartLimitIntervalSec=300
StartLimitBurst=10
User=root
PrivateTmp=true

[Install]
WantedBy=multi-user.target
