# backend/systemd/cmp-cert-expiry.timer
[Unit]
Description=Run CMP certificate expiry check daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=cmp-cert-expiry.service

[Install]
WantedBy=timers.target
