# backend/systemd/cmp-telegram-bot.service
[Unit]
Description=CMP Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/cmp/backend
Environment=NODE_ENV=production
EnvironmentFile=/srv/cmp/backend/.env
ExecStart=/usr/bin/node telegram_bot.js
Restart=on-failure
RestartSec=5s
User=root
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
