# backend/systemd/cmp-matview-refresh.service
[Unit]
Description=CMP Manual Matview Refresh Trigger
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/cmp/backend
Environment=NODE_ENV=production
EnvironmentFile=/srv/cmp/backend/.env
ExecStart=/usr/bin/node scripts/refresh_matview_once.js
User=root

[Install]
WantedBy=multi-user.target